# operation_support/aquaclean-bridge.service
# systemd service unit for aquaclean-bridge
#
# Before installing, replace the two placeholders:
#   YOUR_USER   — the OS user that owns the venv (e.g. jens, pi, aquaclean)
#   /path/to/venv — absolute path to the Python venv (e.g. /home/jens/venv)
#
#   sed -e 's|YOUR_USER|jens|g' \
#       -e 's|/path/to/venv|/home/jens/venv|g' \
#       aquaclean-bridge.service | sudo tee /etc/systemd/system/aquaclean-bridge.service
#
# Then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable aquaclean-bridge
#   sudo systemctl start aquaclean-bridge
#
# Usage:
#   sudo systemctl status aquaclean-bridge
#   sudo systemctl restart aquaclean-bridge
#   tail -f /var/log/aquaclean/aquaclean.log

[Unit]
Description=Geberit AquaClean bridge (BLE → MQTT / REST API)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER

# The '+' prefix runs these as root regardless of User=, ensuring the log
# directory exists and is owned by YOUR_USER before stdout is redirected.
ExecStartPre=+/bin/mkdir -p /var/log/aquaclean
ExecStartPre=+/bin/chown YOUR_USER /var/log/aquaclean

ExecStart=/path/to/venv/bin/aquaclean-bridge --mode api

Restart=on-failure
RestartSec=10

# Log both stdout (JSON CLI results) and stderr (Python logging) to the same file.
StandardOutput=append:/var/log/aquaclean/aquaclean.log
StandardError=append:/var/log/aquaclean/aquaclean.log

[Install]
WantedBy=multi-user.target
